# roles/base/templates/maintenance-reboot.timer.j2
[Unit]
Description=Reboot when required during maintenance window

[Timer]
OnCalendar={{ maintenance_reboot_schedule }}

[Install]
WantedBy=timers.target
