# roles/base/templates/maintenance-reboot.service.j2
[Unit]
Description=Reboot for maintenance purposes

[Service]
Type=notify
ExecStart=/usr/local/sbin/maintenance.sh
